# target/linux/imx6ull/base-files/etc/uci-defaults/99-stream-wifi-setup
#!/bin/sh
#
# Copyright (C) 2013-2015 OpenWrt.org
#

. /lib/functions/uci-defaults.sh

board=$(board_name)

board_config_update

case "$board" in
video_stream_wifi|\
audio_stream_wifi|\
flexcan_wifi|\
lorawan_gateway_wifi)
	touch /etc/config/wireless;
	uci set wireless.default_radio0=wifi-iface
	uci set wireless.@wifi-iface[0].mode="ap";
	uci set wireless.@wifi-iface[0].ifname="wlan0";
	uci set wireless.@wifi-iface[0].network="lan";
#	touch /etc/config/network;
#	uci set network.wlan='interface';
#	uci set network.wlan.proto='dhcp';
	;;
*)
	;;
esac

board_config_flush

exit 0
